# network/ansible/templates/backup-services.timer.j2
[Unit]
Description=Nightly backup of network services to NAS

[Timer]
OnCalendar=*-*-* {{ backup_time }}
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
