# compute/smart_resume_lxc/templates/smart-resume-api.service.j2
[Unit]
Description=Smart Resume API
After=network.target

[Service]
Type=exec
User=www-data
WorkingDirectory={{ app_remote_path }}/api
EnvironmentFile={{ app_remote_path }}/api/.env
Environment=PYTHONPATH={{ app_remote_path }}/api/src
ExecStart={{ app_remote_path }}/api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port {{ api_port }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
